# FumeBot-Rpi3B/fumebot.service
[unit]
Description=FumeBot Main Service
After=multi-user.target

[Service]
Type=idle
Restart=always
WorkingDirectory=/home/pi/FumeBot/
ExecStart=/usr/bin/python2 /home/pi/FumeBot/FumeBotMain.py

[Install]
WantedBy=multi-user.target

